# mac-cron/com.jv.hermes-telegram-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jv.hermes-telegram-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jv/Documents/GitHub/jv-hermes/mac-cron/hermes-telegram-sync.sh</string>
    </array>

    <!--
      Run every 3 hours. launchd fires whenever the Mac is awake at the
      scheduled moment; missed firings replay on wake thanks to
      StartCalendarInterval behavior (at next scheduled tick).
    -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>15</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/jv/Library/Logs/jv-hermes/telegram-sync.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jv/Library/Logs/jv-hermes/telegram-sync.stderr.log</string>

    <!-- Ensure PATH contains common locations for the `railway` CLI. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
